# 运营中枢/scripts/com.karuo.docker-stacks-up.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.karuo.docker-stacks-up</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>sleep 45 &amp;&amp; "/Users/karuo/Documents/个人/卡若AI/运营中枢/scripts/docker_stacks_bootstrap.sh"</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/karuo-docker-stacks-up.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/karuo-docker-stacks-up.err</string>
</dict>
</plist>
