# 01_卡资（金）/金仓_存储备份/Gitea管理/脚本/com.karuo.ai.sync-gitea.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.karuo.ai.sync-gitea</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/karuo/Documents/个人/卡若AI/01_卡资（金）/金仓_存储备份/Gitea管理/脚本/自动同步.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>WorkingDirectory</key>
  <string>/Users/karuo/Documents/个人/卡若AI</string>
  <key>StandardOutPath</key>
  <string>/Users/karuo/Documents/个人/卡若AI/01_卡资（金）/金仓_存储备份/Gitea管理/脚本/sync_launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/karuo/Documents/个人/卡若AI/01_卡资（金）/金仓_存储备份/Gitea管理/脚本/sync_launchd_stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
