# 运营中枢/工作台/脚本/com.openclaw.gateway.longmao.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.gateway.longmao</string>
  <key>Comment</key>
  <string>OpenClaw 龙虾网关（飞书龙猫）— 阿猫 Mac；须 Node &gt;=22.16，故使用 Homebrew Node，勿用 ~/.local/node22</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/opt/node/bin/node</string>
    <string>/Users/kr/.local/share/npm-global/lib/node_modules/openclaw/dist/index.js</string>
    <string>gateway</string>
    <string>--port</string>
    <string>18789</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/kr/.openclaw</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>
  <key>StandardOutPath</key>
  <string>/Users/kr/.openclaw/launchd-gateway-longmao.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kr/.openclaw/launchd-gateway-longmao.err.log</string>
</dict>
</plist>
