# scripts/com.karuo.ai.push-github.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.karuo.ai.push-github</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/karuo/Documents/个人/卡若AI/scripts/push-to-github.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>14400</integer>
  <key>WorkingDirectory</key>
  <string>/Users/karuo/Documents/个人/卡若AI</string>
  <key>StandardOutPath</key>
  <string>/tmp/karuo-ai-push-github.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/karuo-ai-push-github.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
