# 01_卡资（金）/金仓_存储备份/群晖NAS管理/scripts/com.karuo.mount_diskstation_1tb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.karuo.mount_diskstation_1tb</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/karuo/Documents/个人/卡若AI/01_卡资（金）/金仓_存储备份/群晖NAS管理/scripts/mount_diskstation_1tb.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/mount_diskstation_1tb.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mount_diskstation_1tb.err</string>
</dict>
</plist>
