# 01_卡资（金）/金仓_存储备份/群晖NAS管理/scripts/com.karuo.time_machine_home_auto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.karuo.time_machine_home_auto</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/karuo/Documents/个人/卡若AI/01_卡资（金）/金仓_存储备份/群晖NAS管理/scripts/time_machine_home_auto.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>/tmp/time_machine_home_auto.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/time_machine_home_auto.err</string>
</dict>
</plist>
